# templates/systemd/agent-homelab-maintenance.service
[Unit]
Description=Validate Agent Homelab inventory and report drift
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=%h/.config/agent-homelab
ExecStart=%h/.local/bin/agent-homelab maintenance --inventory homelab.yaml
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
