# templates/systemd/agent-homelab-maintenance.timer
[Unit]
Description=Run Agent Homelab maintenance checks

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
